# daemon/pocket-claude.service
[Unit]
Description=pocket-claude — Claude Code always-on daemon
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=forking
User=claude
Group=claude

# Fixed socket path outside PrivateTmp so watchdog + switch can reach it
Environment=TMUX_TMPDIR=/home/claude/.pocket-claude/tmux

# Start Claude in a detached tmux session
ExecStart=/home/claude/.pocket-claude/start.sh
ExecStop=/usr/bin/tmux kill-session -t pocket-claude
PIDFile=/home/claude/.pocket-claude/daemon.pid

# Restart on any exit so session switches (which kill claude) auto-restart
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/claude
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pocket-claude

[Install]
WantedBy=multi-user.target
